# etc/udev/rules.d/99-automount.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[b-z][0-9]", RUN+="/bin/systemctl start usb-handler.service"


